# config/installer/firstboot.service
[Unit]
Description=DayShield First-Boot Initialization
Documentation=https://dayshield.local/docs
After=local-fs.target network-online.target
Wants=network-online.target
Requires=local-fs.target
ConditionPathExists=/etc/dayshield/.firstboot

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/dayshield-installer/firstboot-run.sh
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
